# deploy/alt-bitnodes-archive.timer
[Unit]
Description=Daily snapshot archive run

[Timer]
# After midnight UTC so yesterday's last snapshot exists; jitter to be a
# polite tenant on the shared host.
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
