# deploy/bitnodes.service
[Unit]
Description=Bitnodes crawler stack (crawl + ping + resolve + export + seeder + cache_inv)
After=network-online.target redis-server.service
# tcpdump-pcap.service is deliberately NOT in Wants=. The sniffer's I/O
# contention makes snapshot counts oscillate (postmortem 2026-05-13).
# A Wants= dependency here would pull it in on every crawler restart even
# while `systemctl disable`d — `disable` only removes boot symlinks, not
# Wants=. pcap capture is opt-in: start tcpdump-pcap.service manually if
# ever needed. See change fix-tcpdump-revival.
Wants=network-online.target
Requires=redis-server.service

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__CRAWLER_DIR__
ExecStart=__CRAWLER_DIR__/run-bitnodes.sh
Restart=on-failure
RestartSec=10
KillMode=control-group
TimeoutStopSec=30

LimitNOFILE=65535
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=__CRAWLER_DIR__

[Install]
WantedBy=multi-user.target
