# deploy/tcpdump-pcap.service
[Unit]
Description=tcpdump PCAP producer for bitnodes cache_inv (Bitcoin magic-byte filter)
After=network-online.target bitnodes.service
Wants=network-online.target
PartOf=bitnodes.service

[Service]
Type=simple
ExecStart=__CRAWLER_DIR__/run-tcpdump.sh
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=10

NoNewPrivileges=no
ProtectHome=read-only
ReadWritePaths=__CRAWLER_DIR__/data/pcap

[Install]
WantedBy=multi-user.target
